# .config/systemd/user/colors-wallpaper.service
[Unit]
Description=Changes system colors and wallpaper every hour.

[Service]
ExecStart=/usr/local/bin/colors-wallpaper
Environment=DISPLAY=:0
