# .config/systemd/user/hourly.timer
[Unit]
Description=Hourly timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=hourly.target

[Install]
WantedBy=hourly.target
