# .config/systemd/user/hourly.target
[Unit]
Description=Hourly timer
StopWhenUnneeded=yes
